# conf/supervisor.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /deploy/uwsgi.ini --py-auto-reload 1
